# unicorn-loadbalancer/rootfs/etc/services.d/unicorn-loadbalancer/run
#!/usr/bin/with-contenv sh
npm start --prefix=/opt/UnicornLoadBalancer
